# wkd_sync/fujilinux-keyring-wkd-sync.timer
[Unit]
Description=Refresh existing PGP keys of fujilinux-keyring regularly

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1week
